# .service
[Unit]
Description=PythonTestApplication
After=network.target

[Service]
User=test_user
#Group=users
Group=www-data
PIDFile = /run/PythonTestApplication/PytSample.pid
WorkingDirectory=/home/test_user/app
Environment=POSTGRESQL_URL=postgresql://test_user:DMg00LhU@127.0.0.1/app
Environment=FLASK_APP=app.py
ExecStart = /home/test_user/app/usr/local/bin/gunicorn app:api -b 0.0.0.0:5000 --pid /run/PythonTestApplication/PytSample.pid

[Install]
WantedBy=multi-user.target
