# .service
[Unit]
Description=PythonTestApplication
After=network.target

[Service]
User=test_user
Group=users
PIDFile = /run/PythonTestApplication/PytSample.pid
WorkingDirectory=/home/user/1
Environment="PATH=/home/{{ ansible_ssh_user }}/{{ app_name }}/env/bin"
ExecStart = /usr/bin/env gunicorn sample:app -b 0.0.0.0:5000 --pid /run/PythonTestApplication/PytSample.pid
ExecStart=/home/{{ ansible_ssh_user }}/{{ app_name }}/env/bin/gunicorn --workers 3 --bind unix:{{ app_name }}.sock -m 007 wsgi:app

[Install]
WantedBy=multi-user.target


https://bartsimons.me/gunicorn-as-a-systemd-service/
